# crypto-news-bot.timer
# Systemd timer file for Crypto News Bot
# Runs the bot every 5 minutes
#
# Installation instructions:
# 1. Copy crypto-news-bot.service to /etc/systemd/system/
# 2. Copy this file to /etc/systemd/system/crypto-news-bot.timer
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable crypto-news-bot.timer
# 5. Run: sudo systemctl start crypto-news-bot.timer
# 6. Check status: sudo systemctl status crypto-news-bot.timer
# 7. List next scheduled runs: systemctl list-timers crypto-news-bot.timer
#
# Viewing logs (all timer runs are automatically logged):
# - View recent logs: sudo journalctl -u crypto-news-bot.service -n 50
# - Follow logs: sudo journalctl -u crypto-news-bot.service -f
# - Today's logs: sudo journalctl -u crypto-news-bot.service --since today
# - Last hour: sudo journalctl -u crypto-news-bot.service --since "1 hour ago"
# - Errors only: sudo journalctl -u crypto-news-bot.service -p err
# - Export logs: sudo journalctl -u crypto-news-bot.service > logs.txt
#
# Helper scripts:
# - ./view_logs.sh     - View logs dashboard
# - ./analyze_logs.sh  - Analyze timer execution history

[Unit]
Description=Crypto News Bot Timer (Every 5 Minutes)
Requires=crypto-news-bot.service

[Timer]
# Run every 5 minutes
OnCalendar=*:0/5
# Alternative schedules:
# OnCalendar=hourly              # Every hour on the hour
# OnCalendar=*:0/10              # Every 10 minutes
# OnCalendar=*:0/15              # Every 15 minutes
# OnCalendar=*-*-* *:15:00       # Every hour at :15 past
# If the system was off when the timer was supposed to run, run it as soon as possible
Persistent=true
# Small randomization to avoid exact simultaneous starts (in seconds)
RandomizedDelaySec=30

[Install]
WantedBy=timers.target
