# crypto-news-bot.timer
# Systemd timer file for Crypto News Bot
# Runs the bot every hour
#
# Installation instructions:
# 1. Copy crypto-news-bot.service to /etc/systemd/system/
# 2. Copy this file to /etc/systemd/system/crypto-news-bot.timer
# 3. Run: sudo systemctl daemon-reload
# 4. Run: sudo systemctl enable crypto-news-bot.timer
# 5. Run: sudo systemctl start crypto-news-bot.timer
# 6. Check status: sudo systemctl status crypto-news-bot.timer
# 7. View logs: sudo journalctl -u crypto-news-bot.service
# 8. List next scheduled runs: systemctl list-timers crypto-news-bot.timer

[Unit]
Description=Crypto News Bot Hourly Timer
Requires=crypto-news-bot.service

[Timer]
# Run every hour on the hour (00:00, 01:00, 02:00, etc.)
OnCalendar=hourly
# Alternative: Run every hour at a specific minute (e.g., at :15 past each hour)
# OnCalendar=*-*-* *:15:00
# If the system was off when the timer was supposed to run, run it as soon as possible
Persistent=true
# Small randomization to avoid exact simultaneous starts
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
